# phone-home-indicator.timer
[Unit]
Description=Trigger Phone home indicator

[Timer]
OnUnitInactiveSec=3600

[Install]
WantedBy=timers.target
