# phone-home-indicator.timer
[Unit]
Description=Trigger Phone home indicator

[Timer]
OnUnitInactiveSec=15s

[Install]
WantedBy=timers.target
